# System_files/etc/systemd/system/mavlink.router.service
[Unit]
Description=MAVLink Router
After=network.target wifibroadcast@gs.service
Wants=wifibroadcast@gs.service

[Service]
ExecStart=/usr/bin/mavlink-routerd -c /etc/mavlink-router/main.conf
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
